# ops/docker/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# NOTE:
# - These jobs assume /app/.env exists (mount it) so python-dotenv can load secrets.
# - Logs go to Data/logs/*.cron.log (mount ./Data for persistence).

# Hourly forecast pulse (append-only) at :05
05 * * * * cd /app && /bin/bash /app/scripts/run_hourly_pulse.sh >> /app/Data/logs/hourly_pulse.cron.log 2>&1

# Daily trade run for tomorrow's markets at 10:00pm local container time (America/New_York)
00 22 * * * cd /app && /bin/bash /app/scripts/run_trade.sh >> /app/Data/logs/trade.cron.log 2>&1

# Nightly calibration for yesterday (default 2:15am local container time)
15 2 * * * cd /app && /bin/bash /app/scripts/run_calibrate.sh >> /app/Data/logs/calibrate.cron.log 2>&1

# Nightly settlement + metrics rollups for yesterday (default 3:15am local container time)
15 3 * * * cd /app && /bin/bash /app/scripts/run_settle.sh >> /app/Data/logs/settle.cron.log 2>&1

