# ops/docker/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# NOTE:
# - These jobs assume /app/.env exists (mount it) so python-dotenv can load secrets.
# - Logs go to Data/logs/*.cron.log (mount ./Data for persistence).

# Intraday forecast pulse (append-only) at 09:00, 15:00, 21:00 local container time (America/New_York)
# Stream to Docker logs AND append to file
00 09 * * * cd /app && /bin/bash /app/scripts/run_intraday_pulse.sh 2>&1 | tee -a /app/Data/logs/intraday_pulse.cron.log
00 15 * * * cd /app && /bin/bash /app/scripts/run_intraday_pulse.sh 2>&1 | tee -a /app/Data/logs/intraday_pulse.cron.log
00 21 * * * cd /app && /bin/bash /app/scripts/run_intraday_pulse.sh 2>&1 | tee -a /app/Data/logs/intraday_pulse.cron.log

# Daily trade run for tomorrow's markets at 10:30pm local container time (America/New_York)
# (This run performs the final 22:00 forecast fetch + trades if gates pass.)
30 22 * * * cd /app && /bin/bash /app/scripts/run_trade.sh 2>&1 | tee -a /app/Data/logs/trade.cron.log

# Nightly calibration for yesterday (default 2:15am local container time)
15 2 * * * cd /app && /bin/bash /app/scripts/run_calibrate.sh 2>&1 | tee -a /app/Data/logs/calibrate.cron.log

# Nightly settlement + metrics rollups for yesterday (default 3:15am local container time)
15 3 * * * cd /app && /bin/bash /app/scripts/run_settle.sh 2>&1 | tee -a /app/Data/logs/settle.cron.log

# Retry window (fault tolerance): some CLIs are published later (e.g., MIA can be ~4am local,
# AUS is CST, etc.). Re-run later in the morning to pick up late-issued reports.
30 6 * * * cd /app && /bin/bash /app/scripts/run_calibrate.sh 2>&1 | tee -a /app/Data/logs/calibrate.cron.log
15 7 * * * cd /app && /bin/bash /app/scripts/run_settle.sh 2>&1 | tee -a /app/Data/logs/settle.cron.log

