# ops/docker/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# NOTE:
# - These jobs assume /app/.env exists (mount it) so python-dotenv can load secrets.
# - Logs go to Data/logs/*.cron.log (mount ./Data for persistence).

# Intraday forecast pulse: 41 runs/day (within free-tier limits; see documentation/forecast_api_limits.md).
# Binding limit: Google 10,000/month → 41 runs/day (8 calls/run × 41 × 30 ≈ 9,840). Tomorrow.io 25/hr, 500/day OK.
# Every hour at :00 (24/day) + at :30 for hours 0-16 (17/day) = 41 runs/day.
0 * * * * cd /app && /bin/bash /app/scripts/run_intraday_pulse.sh 2>&1 | tee -a /app/Data/logs/intraday_pulse.cron.log
30 0-16 * * * cd /app && /bin/bash /app/scripts/run_intraday_pulse.sh 2>&1 | tee -a /app/Data/logs/intraday_pulse.cron.log

# Daily trade run for **today's** markets. Container TZ = America/New_York.
# - 07:00 ET: early-morning forecast refresh; kalshi_trader 13:00 gate skips all (use --skip-13h-gate for one-shot trade-all).
# - 13:00 ET: run_trade runs intraday_pulse (fresh forecasts) then kalshi_trader; ny/fl are 13:00 local → trade.
# - 14:00 ET: run_trade runs intraday_pulse (fresh forecasts) then kalshi_trader; il/tx are 13:00 CT → trade.
00 7 * * * cd /app && /bin/bash /app/scripts/run_trade.sh 2>&1 | tee -a /app/Data/logs/trade.cron.log
0 13 * * * cd /app && /bin/bash /app/scripts/run_trade.sh 2>&1 | tee -a /app/Data/logs/trade.cron.log
0 14 * * * cd /app && /bin/bash /app/scripts/run_trade.sh 2>&1 | tee -a /app/Data/logs/trade.cron.log

# Nightly calibration for yesterday (default 2:15am local container time)
15 2 * * * cd /app && /bin/bash /app/scripts/run_calibrate.sh 2>&1 | tee -a /app/Data/logs/calibrate.cron.log

# Nightly settlement + metrics rollups for yesterday (default 3:15am local container time)
15 3 * * * cd /app && /bin/bash /app/scripts/run_settle.sh 2>&1 | tee -a /app/Data/logs/settle.cron.log

# Retry window (fault tolerance): some CLIs are published later (e.g., MIA can be ~4am local,
# AUS is CST, etc.). Re-run later in the morning to pick up late-issued reports.
30 6 * * * cd /app && /bin/bash /app/scripts/run_calibrate.sh 2>&1 | tee -a /app/Data/logs/calibrate.cron.log
15 7 * * * cd /app && /bin/bash /app/scripts/run_settle.sh 2>&1 | tee -a /app/Data/logs/settle.cron.log

